# contrib/systemd/quicport.service
[Unit]
Description=quicport - QUIC-based port forwarding / tunneling server
Documentation=https://github.com/hakadoriya/quicport
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Create a group and user for the service
# $ sudo groupadd --system quicport
# $ sudo useradd --system --gid quicport --shell /usr/sbin/nologin --create-home --home-dir /var/lib/quicport --comment "quicport service user" quicport
User=quicport
Group=quicport

# NOTE: The binary path
ExecStart=/usr/local/bin/quicport --log-format json server
ExecReload=/usr/local/bin/quicport ctl graceful-restart
KillMode=process

# NOTE: The environment variable file (authentication key, etc.)
EnvironmentFile=-/etc/quicport/quicport.env

# NOTE: Restart settings
Restart=on-failure
RestartSec=5

# NOTE: Log settings
# Create a log directory
# $ sudo mkdir -m 755 -p /var/log/quicport
# $ sudo chown quicport:quicport /var/log/quicport
StandardOutput=append:/var/log/quicport/quicport.log
StandardError=inherit
SyslogIdentifier=quicport

# NOTE: File descriptor limit
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
